# audio_logger.service
[Unit]
Description=Audio Logger Service
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=audiologger
Group=audiologger
WorkingDirectory=/opt/audio_logger
ExecStart=/opt/audio_logger/.venv/bin/python /opt/audio_logger/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
UMask=0027

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
